# restic-check@.timer
[Unit]
Description=Check Restic backup daily

[Timer]
OnCalendar=daily
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
