# restic-prune@.timer
[Unit]
Description=Prune Restic backup weekly

[Timer]
OnCalendar=daily
RandomizedDelaySec=6h
 
[Install]
WantedBy=timers.target
